# deploy/oci-resource-extractor-web.service.example
[Unit]
Description=OCI Resource Extractor Web
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=oci-extractor
Group=oci-extractor
WorkingDirectory=/opt/OCI-ResourceExtractor
Environment=PATH=/opt/OCI-ResourceExtractor/venv/bin
Environment=OCI_API_MIN_INTERVAL_MS=0
# Uncomment and set private values outside Git when web login is required.
# EnvironmentFile=-/etc/oci-resource-extractor/web.env
# Environment=OCI_EXTRACTOR_AUTH_ENABLED=true
# Environment=OCI_EXTRACTOR_USERNAME=admin
# Environment=OCI_EXTRACTOR_PASSWORD_HASH=pbkdf2_sha256$260000$...
# Environment=OCI_EXTRACTOR_SESSION_SECRET=replace-with-private-random-secret
# Environment=OCI_EXTRACTOR_SESSION_TTL_MINUTES=480
# Environment=OCI_EXTRACTOR_COOKIE_SECURE=false
ExecStart=/opt/OCI-ResourceExtractor/venv/bin/uvicorn web_app:app --host 127.0.0.1 --port 8088
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
