# deploy/oci-resource-extractor-web.service.example
[Unit]
Description=OCI Resource Extractor Web
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=oci-extractor
Group=oci-extractor
WorkingDirectory=/opt/OCI-ResourceExtractor
Environment=PATH=/opt/OCI-ResourceExtractor/venv/bin
Environment=OCI_API_MIN_INTERVAL_MS=0
ExecStart=/opt/OCI-ResourceExtractor/venv/bin/uvicorn web_app:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
